# deploy/quadlet/cadets-web.container
[Unit]
Description=763 Cadets web gateway (publishes one port, routes /api + serves SPA)
Wants=cadets-backend.service cadets-frontend.service
After=cadets-backend.service cadets-frontend.service

[Container]
ContainerName=cadets-web
# Built from caddy.Containerfile (Caddyfile baked in -> no fragile config mount).
Image=cadets-web.build
Network=cadets.network
PublishPort=8088:8088/tcp

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
