# deploy/quadlet/cadets-web.container
[Unit]
Description=763 Cadets web gateway (publishes one port, routes /api + serves SPA)
Wants=cadets-backend.service cadets-frontend.service
After=cadets-backend.service cadets-frontend.service

[Container]
ContainerName=cadets-web
Image=docker.io/library/caddy:2-alpine
Network=cadets.network
PublishPort=8088:8088/tcp
# CHANGE the host path to where you cloned the repo.
Volume=/home/alex/cadet763/deploy/quadlet/Caddyfile:/etc/caddy/Caddyfile:ro,Z
AutoUpdate=registry

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
