# deploy/quadlet/cadets-frontend.container
[Unit]
Description=763 Cadets frontend (nginx static SPA)
After=local-fs.target

[Container]
ContainerName=cadets-frontend
# Built by CI and pulled from GHCR.
Image=ghcr.io/alexb715/cadet763-frontend:latest
Pull=newer
AutoUpdate=registry
Network=cadets.network

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
